# olpc-switchd.upstart
#

description     "XO hw switch handler"
author          "Paul Fox"

start on starting olpc-kbdshim

console output

exec /usr/sbin/olpc-switchd -f -l -F /var/run/powerevents

respawn
