# olpc-switchd.upstart
#

description     "XO hw switch handler"
author          "Paul Fox"

start on stopping rc5

console output

# poll for power events every 10 seconds,
# report periodic timer at least every 60
exec /usr/sbin/olpc-switchd -f -l -p 10 -t 6 -F /var/run/powerevents

respawn
