# powerd.upstart
#
description     "XO experimental power daemon "
author          "Paul Fox"

start on starting olpc-switchd

console output

exec /usr/sbin/powerd

respawn
