# powerd.upstart
#
description     "XO experimental power daemon "
author          "Paul Fox"

start on stopping rc RUNLEVEL=5

stop on runlevel [016]

console output

exec /usr/sbin/powerd

respawn
