# powerd.upstart
#
description     "XO experimental power daemon "
author          "Paul Fox"

start on stopping rc5

console output

exec /usr/sbin/powerd

respawn
