# powerd.service
[Unit]
Description=OLPC XO power daemon

[Service]
ExecStart=/usr/sbin/powerd

[Install]
WantedBy=multi-user.target
